# systemd/security-audit.service
[Unit]
Description=Daily Security Audit

[Service]
Type=oneshot
ExecStart=/usr/local/bin/security-audit.sh
